# ops/launchd/co.getforevermore.autopilot.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.getforevermore.autopilot.telegram</string>
	<!-- The Telegram control channel (WAVE2 Phase 1). A long-lived long-poll +
	     scan daemon: a projection of the DB into the owner's phone (review cards
	     with playable media, tick summaries, failure/liveness/spend alerts) and
	     the inbound steering surface (/new, /rule, /pause, freeform notes).
	     Decisions flow through the SAME decide()/store paths as the Station, so
	     the two never disagree. Own lockfile (state/telegram.lock) prevents a
	     second poller. KeepAlive restarts it on crash and at boot.

	     Secrets are NOT in this file (it's committed). Put them in
	     ~/.config/autopilot/telegram.env (chmod 600), which the shell sources:
	         TELEGRAM_ENABLED=true
	         TELEGRAM_BOT_TOKEN=123456:AA...
	         TELEGRAM_CHAT_ID=900900900
	     On the VPS (Phase 2) this becomes a systemd unit with EnvironmentFile. -->
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>cd /Users/redonemini/Desktop/Code/Personal/forevermore-autopilot &amp;&amp; [ -f ~/.config/autopilot/telegram.env ] &amp;&amp; set -a &amp;&amp; . ~/.config/autopilot/telegram.env &amp;&amp; set +a; node bin/autopilot.mjs telegram >> logs/launchd-telegram.log 2>&amp;1</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>30</integer>
	<!-- TimeZone: Europe/Tirane (system TZ) -->
</dict>
</plist>
