# ops/launchd/co.getforevermore.autopilot.reflect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.getforevermore.autopilot.reflect</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>cd /Users/redonemini/Desktop/Code/Personal/forevermore-autopilot &amp;&amp; node bin/autopilot.mjs run reflect >> logs/launchd-reflect.log 2>&amp;1</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>23</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<!-- TimeZone: Europe/Tirane (system TZ) -->
</dict>
</plist>
